# deploy/supervisor.conf
; BRICKStack Supervisor Configuration
; Install: apt install supervisor (or pip install supervisor)
; Place in: /etc/supervisor/conf.d/brickstack.conf

[program:brickstack-backend]
command=python3 -m uvicorn backend.main_secure:app --host 0.0.0.0 --port 8000 --workers 2
directory=/root/BRICKStack
autostart=true
autorestart=true
startretries=3
user=brickstack
stderr_logfile=/var/log/brickstack/backend.err.log
stdout_logfile=/var/log/brickstack/backend.out.log
environment=PATH="/usr/bin:/bin:/usr/local/bin",HOME="/root/BRICKStack"

[program:brickstack-telegram]
command=python3 telegram_bot.py
directory=/root/BRICKStack
autostart=true
autorestart=true
startretries=3
user=brickstack
stderr_logfile=/var/log/brickstack/telegram.err.log
stdout_logfile=/var/log/brickstack/telegram.out.log
environment=PATH="/usr/bin:/bin:/usr/local/bin",HOME="/root/BRICKStack"
